# packaging/app-antivirus.cron
2 * * * * root /usr/sbin/freshclam-update >/dev/null 2>&1
